# roles/mail/templates/mailcatcher.conf.j2
start on startup
exec /usr/local/bin/mailcatcher --foreground --ip 0.0.0.0